# app-admin/jailkit/files/jk_socketd.initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

description="A daemon to create a rate-limited /dev/log socket inside a chroot"

command="/usr/sbin/jk_socketd"
pidfile="/var/run/${SVCNAME}.pid"
command_args="--pidfile=${pidfile}"

required_files="/etc/jailkit/jk_socketd.ini"

depend() {
	use logger
}
